# firejail@.service
[Service]
ExecStart=/usr/bin/firejail --socket-activation
